# optscale_single/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
loglevel=warn
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/optscale

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ============================================================================
# INFRASTRUCTURE — priority 5, start first
# ============================================================================

[program:etcd]
command=etcd
    --enable-v2=true
    --listen-client-urls=http://0.0.0.0:2379
    --advertise-client-urls=http://0.0.0.0:2379
    --data-dir=/default.etcd
    --log-level=warn
priority=5
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/var/log/optscale/etcd.log
stderr_logfile=/var/log/optscale/etcd.log
environment=ETCD_HEARTBEAT_INTERVAL="500",ETCD_ELECTION_TIMEOUT="5000"

[program:mariadb]
command=/usr/sbin/mariadbd
    --user=mysql
    --datadir=/var/lib/mysql
    --max_connections=128
    --innodb_buffer_pool_size=128M
    --innodb_file_per_table=1
    --wait_timeout=28800
    --explicit_defaults_for_timestamp=1
priority=5
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/optscale/mariadb.log
stderr_logfile=/var/log/optscale/mariadb.log

[program:mongod]
command=/usr/bin/mongod
    --replSet=mongo
    --auth
    --keyFile=/data/configdb/key.txt
    --bind_ip_all
    --dbpath=/data/db
    --port=27017
    --storageEngine=wiredTiger
    --wiredTigerCacheSizeGB=0.25
    --logappend
user=mongodb
priority=5
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/optscale/mongod.log
stderr_logfile=/var/log/optscale/mongod.log

[program:redis]
command=/usr/bin/redis-server
    --save ""
    --maxmemory 64mb
    --maxmemory-policy allkeys-lru
    --loglevel warning
priority=5
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/var/log/optscale/redis.log
stderr_logfile=/var/log/optscale/redis.log

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
priority=5
autorestart=true
startretries=10
startsecs=15
stdout_logfile=/var/log/optscale/rabbitmq.log
stderr_logfile=/var/log/optscale/rabbitmq.log
environment=HOME="/var/lib/rabbitmq",RABBITMQ_NODENAME="rabbit@localhost",RABBITMQ_ERLANG_COOKIE="dWkhXehOApH1OHshvByfpfTf3Cd8wwkADei7mWd5JF"

[program:clickhouse]
command=/usr/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml
user=clickhouse
priority=5
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/optscale/clickhouse.log
stderr_logfile=/var/log/optscale/clickhouse.log

[program:influxdb]
command=/usr/bin/influxd -config /etc/influxdb/influxdb.conf
priority=5
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/optscale/influxdb.log
stderr_logfile=/var/log/optscale/influxdb.log

[program:minio]
command=minio server /opt/minio-data --address :9000
priority=5
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/var/log/optscale/minio.log
stderr_logfile=/var/log/optscale/minio.log
environment=MINIO_ACCESS_KEY="optscale-minio",MINIO_SECRET_KEY="secret_password"

# ============================================================================
# CONFIGURATOR — one-shot, runs after infra is healthy (priority 50)
# Has internal retry logic (300 attempts × 500ms) so it waits for DBs.
# ============================================================================

[program:configurator]
command=python3 /usr/src/app/docker_images/configurator/configurator.py /config/etcd-config
directory=/usr/src/app
priority=50
autorestart=false
startretries=0
startsecs=3
exitcodes=0
stdout_logfile=/var/log/optscale/configurator.log
stderr_logfile=/var/log/optscale/configurator.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

# ============================================================================
# CORE API SERVICES — priority 100
# All have internal retry logic for etcd/DB connection.
# ============================================================================

[program:auth]
command=uv --project auth run python -u /usr/src/app/auth/auth_server/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/auth.log
stderr_logfile=/var/log/optscale/auth.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",
    GOOGLE_OAUTH_CLIENT_ID="",GOOGLE_OAUTH_CLIENT_SECRET="",
    MICROSOFT_OAUTH_CLIENT_ID=""

[program:restapi]
command=uv --project rest_api run python -u /usr/src/app/rest_api/rest_api_server/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/restapi.log
stderr_logfile=/var/log/optscale/restapi.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",
    FAKE_CAD_ENABLED="0",
    REQUESTS_CA_BUNDLE="/etc/ssl/certs/ca-certificates.crt"

[program:keeper]
command=uv --project keeper run python -u /usr/src/app/keeper/report_server/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/keeper.log
stderr_logfile=/var/log/optscale/keeper.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:herald-api]
command=uv --project herald run python -u /usr/src/app/herald/run.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/herald-api.log
stderr_logfile=/var/log/optscale/herald-api.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",HERALD_SERVICE="api"

[program:herald-engine]
command=uv --project herald run python -u /usr/src/app/herald/run.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/herald-engine.log
stderr_logfile=/var/log/optscale/herald-engine.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",HERALD_SERVICE="engine"

[program:insider-api]
command=uv --project insider/insider_api run python -u /usr/src/app/insider/insider_api/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/insider-api.log
stderr_logfile=/var/log/optscale/insider-api.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:metroculus-api]
command=uv --project metroculus/metroculus_api run python -u /usr/src/app/metroculus/metroculus_api/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/metroculus-api.log
stderr_logfile=/var/log/optscale/metroculus-api.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:slacker]
command=uv --project slacker run python /src/slacker/slacker_server/server.py
directory=/src
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/slacker.log
stderr_logfile=/var/log/optscale/slacker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",
    SLACKER_PORT="8960",
    SLACK_SIGNING_SECRET="",SLACK_CLIENT_ID="",SLACK_CLIENT_SECRET=""

[program:diproxy]
command=uv --project diproxy run python -u /usr/src/app/diproxy/diproxy/main.py --port=8936
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/diproxy.log
stderr_logfile=/var/log/optscale/diproxy.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:jira-bus]
command=uv --project jira_bus run python /usr/src/app/jira_bus/jira_bus_server/server.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/jira-bus.log
stderr_logfile=/var/log/optscale/jira-bus.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:katara-api]
command=uv --project katara/katara_service run python -u /usr/src/app/katara/katara_service/main.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/katara-api.log
stderr_logfile=/var/log/optscale/katara-api.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",HX_KATARA_ROLE="api"

[program:katara-scheduler]
command=uv --project katara/katara_service run python -u /usr/src/app/katara/katara_service/main.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/katara-scheduler.log
stderr_logfile=/var/log/optscale/katara-scheduler.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",HX_KATARA_ROLE="scheduler"

[program:bi-exporter]
command=uv run --project bi_exporter python -u /usr/src/app/bi_exporter/bumblebi/exporter/main.py
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/bi-exporter.log
stderr_logfile=/var/log/optscale/bi-exporter.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:subspector]
command=uv --project subspector run uvicorn
    subspector.subspector_server.server:make_app
    --factory --host 0.0.0.0 --port 8800
directory=/usr/src/app
priority=100
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/subspector.log
stderr_logfile=/var/log/optscale/subspector.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

# ── ngui (main frontend) ──────────────────────────────────────────────────────
[program:ngui]
command=bash /usr/src/app/prepare-and-run.sh
directory=/usr/src/app
priority=150
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/optscale/ngui.log
stderr_logfile=/var/log/optscale/ngui.log
environment=UI_BUILD_PATH="/usr/src/app/ui",
    PROXY_URL="http://127.0.0.1:80",
    AUTH_ENDPOINT="http://127.0.0.1:8905",
    RESTAPI_ENDPOINT="http://127.0.0.1:8999",
    KEEPER_ENDPOINT="http://127.0.0.1:8973",
    SLACKER_ENDPOINT="http://127.0.0.1:8960",
    VITE_GOOGLE_OAUTH_CLIENT_ID="",
    VITE_MICROSOFT_OAUTH_CLIENT_ID="",
    VITE_KEYCLOAK_URL="",
    VITE_KEYCLOAK_REALM="",
    VITE_KEYCLOAK_CLIENT_ID="",
    VITE_ON_INITIALIZE_ORGANIZATION_SETUP_MODE="automatic",
    VITE_FINOPS_IN_PRACTICE_PORTAL_OVERVIEW="",
    VITE_BILLING_INTEGRATION="",
    VITE_APP_THEME=""

# ── jira_ui (Jira integration frontend) ──────────────────────────────────────
[program:jira-ui]
command=node index.js
directory=/usr/src/app/jira_ui/server
priority=150
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/optscale/jira-ui.log
stderr_logfile=/var/log/optscale/jira-ui.log
environment=PORT="4001",HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",UI_BUILD_PATH="/usr/src/app/jira_ui/ui"

# ============================================================================
# WORKERS — priority 150, long-running queue consumers
# ============================================================================

[program:bumischeduler]
command=uv run --project bumischeduler python -u /usr/src/app/bumischeduler/bumischeduler/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/bumischeduler.log
stderr_logfile=/var/log/optscale/bumischeduler.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:bumiworker]
command=uv run --project bumiworker python -u /usr/src/app/bumiworker/bumiworker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/bumiworker.log
stderr_logfile=/var/log/optscale/bumiworker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",FAKE_CAD_ENABLED="0"

[program:diworker]
command=uv run --project diworker /src/diworker/diworker/main.py
directory=/src
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/diworker.log
stderr_logfile=/var/log/optscale/diworker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",FAKE_CAD_ENABLED="0"

[program:gemini-worker]
command=uv run --project gemini/gemini_worker python -u /usr/src/app/gemini/gemini_worker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/gemini-worker.log
stderr_logfile=/var/log/optscale/gemini-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:insider-worker]
command=uv run --project insider/insider_worker python -u /usr/src/app/insider/insider_worker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/insider-worker.log
stderr_logfile=/var/log/optscale/insider-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:metroculus-worker]
command=uv run --project metroculus/metroculus_worker python -u /usr/src/app/metroculus/metroculus_worker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/metroculus-worker.log
stderr_logfile=/var/log/optscale/metroculus-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:risp-worker]
command=uv run --project risp/risp_worker python -u /usr/src/app/risp/risp_worker/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/risp-worker.log
stderr_logfile=/var/log/optscale/risp-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:trapper-worker]
command=uv run --project trapper/trapper_worker python -u /usr/src/app/trapper/trapper_worker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/trapper-worker.log
stderr_logfile=/var/log/optscale/trapper-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:katara-worker]
command=uv run --project katara/katara_worker python -u /usr/src/app/katara/katara_worker/main.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/katara-worker.log
stderr_logfile=/var/log/optscale/katara-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379"

[program:resource-discovery-worker]
command=python3 docker_images/resource_discovery/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/resource-discovery-worker.log
stderr_logfile=/var/log/optscale/resource-discovery-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",FAKE_CAD_ENABLED="0",PYTHONPATH="/usr/src/app"

[program:resource-observer-worker]
command=python3 docker_images/resource_observer/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/resource-observer-worker.log
stderr_logfile=/var/log/optscale/resource-observer-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:resource-violations-worker]
command=python3 docker_images/resource_violations/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/resource-violations-worker.log
stderr_logfile=/var/log/optscale/resource-violations-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:booking-observer-worker]
command=python3 docker_images/booking_observer/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/booking-observer-worker.log
stderr_logfile=/var/log/optscale/booking-observer-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:organization-violations-worker]
command=python3 docker_images/organization_violations/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/organization-violations-worker.log
stderr_logfile=/var/log/optscale/organization-violations-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:power-schedule-worker]
command=python3 docker_images/power_schedule/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/power-schedule-worker.log
stderr_logfile=/var/log/optscale/power-schedule-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:calendar-observer-worker]
command=python3 docker_images/calendar_observer/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/calendar-observer-worker.log
stderr_logfile=/var/log/optscale/calendar-observer-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:live-demo-generator-worker]
command=python3 docker_images/live_demo_generator/worker.py
directory=/usr/src/app
priority=150
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/live-demo-generator-worker.log
stderr_logfile=/var/log/optscale/live-demo-generator-worker.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

# ============================================================================
# EXECUTORS — priority 200, message-driven action dispatchers
# ============================================================================

[program:herald-executor]
command=python3 docker_images/herald_executor/worker.py
directory=/usr/src/app
priority=200
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/herald-executor.log
stderr_logfile=/var/log/optscale/herald-executor.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:keeper-executor]
command=python3 docker_images/keeper_executor/worker.py
directory=/usr/src/app
priority=200
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/keeper-executor.log
stderr_logfile=/var/log/optscale/keeper-executor.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:slacker-executor]
command=python3 docker_images/slacker_executor/worker.py
directory=/usr/src/app
priority=200
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/slacker-executor.log
stderr_logfile=/var/log/optscale/slacker-executor.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

[program:webhook-executor]
command=python3 docker_images/webhook_executor/worker.py
directory=/usr/src/app
priority=200
autorestart=true
startretries=30
startsecs=5
stdout_logfile=/var/log/optscale/webhook-executor.log
stderr_logfile=/var/log/optscale/webhook-executor.log
environment=HX_ETCD_HOST="127.0.0.1",HX_ETCD_PORT="2379",PYTHONPATH="/usr/src/app"

# ============================================================================
# NGINX — priority 250, starts last
# ============================================================================

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=250
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/var/log/optscale/nginx.log
stderr_logfile=/var/log/optscale/nginx.log
